# install/ipfs.service
[Unit]
Description=IPFS daemon
After=network.target

[Service]
ExecStart=/snap/bin/ipfs daemon

[Install]
WantedBy=multiuser.target
